# deploy/kriteria.service
[Unit]
Description=Kriteria — criteria-based student evaluation
After=network.target

[Service]
WorkingDirectory=/home/ubuntu/Projects/kriteria
ExecStart=/home/ubuntu/Projects/kriteria/bin/kriteria --addr :80
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
